# sixaxis@.service.in
[Unit]
Description=sixaxis helper (%I)

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/sixaxis-helper.sh %I

[Install]
WantedBy=multi-user.target
